# supervisor.conf
[supervisord]
nodaemon=true

[program:zk]
command=/bin/sh /usr/local/bin/zk-init.sh
